# httpd/files/etc/cron.daily/httpd
#!/bin/sh

find /var/log/httpd/ -name "*_log.*" -type f -mtime +30 -exec rm -f {} \;
