# lvs/files/service/keepalived-check/run
#!/bin/sh

PID=/var/run/keepalived_check.pid
[ -e $PID ] && exit
exec /usr/sbin/keepalived --dont-fork --log-detail \
                          --check \
                          -p $PID 2>&1
