# lvs/files/service/keepalived-vrrp/run
#!/bin/sh

PID=/var/run/keepalived_vrrp.pid
[ -e $PID ] && exit
exec /usr/sbin/keepalived --dont-fork --log-detail \
                          --vrrp \
                          -p $PID 2>&1
